# deploy/launchd/com.innsaeit.cockpit-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.innsaeit.cockpit-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>__REPO__/.venv/bin/python</string>
        <string>-m</string>
        <string>pipeline</string>
        <string>--loop</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_API_KEY</key>
        <string>__ANTHROPIC_API_KEY__</string>
        <key>GEMINI_API_KEY</key>
        <string>__GEMINI_API_KEY__</string>
        <key>GROQ_API_KEY</key>
        <string>__GROQ_API_KEY__</string>
        <key>OPENROUTER_API_KEY</key>
        <string>__OPENROUTER_API_KEY__</string>
        <key>OPENAI_API_KEY</key>
        <string>__OPENAI_API_KEY__</string>
        <key>APIFY_TOKEN</key>
        <string>__APIFY_TOKEN__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/cockpit-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/cockpit-watcher.log</string>
</dict>
</plist>
